# files/etc/cron.d/aprs-pigate
* * * * * pi /home/pi/dw-start-cron.sh >/dev/null 2>&1
